# factory/usr/lib/sysmaster/system/hostname-setup.service
[Unit]
Description=hostname setup

[Service]
ExecStart=-/usr/lib/sysmaster/hostname_setup
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
