# factory/usr/lib/sysmaster/system/paths.target
[Unit]
Description=Path Units
Documentation=man sysmaster special
