# factory/usr/lib/sysmaster/system/syslog.target
[Unit]
Description=Initialize syslog
DefaultDependencies=false
