# docs/use/虚机中替代pid1运行/getty.target
[Unit]
Description=Login Prompts
Documentation=man sysmaster secial target
