# docs/use/虚机中替代pid1运行/getty@.service
[Unit]
Description=Getty on %I
ConditionPathExists=/dev/tty0

[Service]
ExecStart=-/sbin/agetty -o '-p -- \\u' --noclear %I
Type=simple
Restart=always
RestartSec=0

[Install]
WantedBy=getty.target
