# docs/use/虚机中替代pid1运行/lvm-activate-openeuler.service
[Unit]
Description=/usr/sbin/lvm vgchange -aay --nohints openeuler

[Service]
RemainAfterExit=true
ExecStart=/usr/sbin/lvm vgchange -aay --nohints openeuler

[Install]
WantedBy=multi-user.target
