# tests/test_units/basic.target
[Unit]
Description=Basic System
Documentation=man sysmaster special
Requires=sysinit.target
