# tests/test_units/multi-user.target
[Unit]
Description=Mutil-user.target
Requires=basic.target
