# tests/test_units/shutdown.target
[Unit]
Description=shutdown target
Documentation=man sysmaster shutdown target
DefaultDependencies=false
