# tests/test_units/test2.service
[Unit]
Description=job test 002
DefaultDependencies=false

[Service]
ExecStart=/bin/sleep 888

[Install]
WantedBy=dbus.service
