# tests/test_units/testsunit.target
[Unit]
Description=Testsunit Description
Requires=Multi-user.target
After=Multi-user.target
Conflicts=rescue.target
