# tests/test_units/test1.service
[Unit]
Description=job test 001
DefaultDependencies=false

[Service]
ExecStart=/bin/sleep 888

[Install]
WantedBy=dbus.service
