# tools/run_with_vm/dbus.socket
[Unit]
Description=D-Bus Socket

[Socket]
ListenStream=/var/run/dbus/system_bus_socket
