# tools/run_with_vm/fstab.service
[Unit]

[Service]
ExecStart="/usr/lib/sysmaster/fstab"
