# tools/run_with_vm/fstab.service
[Unit]
Description=automount fstab device

[Service]
ExecStart=/usr/lib/sysmaster/fstab

[Install]
WantedBy=multi-user.target
