# tools/run_with_vm/loopback-up.service
[Unit]
Description="up loopback device"

[Service]
ExecStart="/usr/sbin/ip link set lo up"
