# tools/run_with_vm/set-hostname.service
[Unit]
Description="set up hostname bydefaul"

[Service]
ExecStart="/usr/lib/sysmaster/set-hostname.sh"

[Install]
WantedBy="multi-user.target"
