# tools/run_with_vm/setip.service
[Unit]
Description="set ip for ssh"

[Service]
ExecStart="/usr/lib/sysmaster/setip.sh"

[Install]
WantedBy="multi-user.target"
