# tools/run_with_busybox/run_image/extra/etc/sysmaster/basic.target
[Unit]
Description=Basic System
Wants=init.service
