# docs/use/兼容busybox模式运行/run_image/extra/etc/sysmaster/init.service
[Unit]
Description=init.service

[Service]
Type=simple
ExecStart=/usr/bin/sh /etc/init.d/rcS
