# docs/use/devmaster替代udev运行/rules.d/90-vconsole.rules
# Each vtcon keeps its own state of fonts.
#
# ACTION=="add", SUBSYSTEM=="vtconsole", KERNEL=="vtcon*", RUN+="/usr/lib/systemd/systemd-vconsole-setup"
