# tests/test_units/tests/after.service
[Unit]
Description=After service
After=fork.service

[Service]
Type=forking
ExecStart=/bin/sleep 5
