# tests/test_units/tests/base.socket
[Unit]
Description=base socket

[Socket]
