# tests/test_units/tests/base.service
[Unit]
Description=base service

[Service]
ExecStart=/bin/sleep 100
