# tests/test_units/tests/before.service
[Unit]
Description=Before service
Before=fork.service

[Service]
Type=forking
ExecStart=/bin/sleep 5
