# tests/test_units/tests/bindsto.service
[Unit]
Description=BindsTo service
BindsTo=base.service

[Service]
ExecStart=/bin/sleep 100
