# tests/test_units/tests/conflicts.service
[Unit]
Description=Conflicts service
Conflicts=base.service

[Service]
ExecStart=/bin/sleep 100
