# tests/test_units/tests/env.service
[Unit]
Description=env service

[Service]
ExecStart=/bin/sleep 100
Environment=VAR1="word1" VAR2="word2" VAR3="\"word 3\"" VAR4="word 4" VAR5="word=5" VAR0="word;0" VAR2="word1"
