# tests/test_units/tests/fork.service
[Unit]
Description=fork service

[Service]
Type=forking
ExecStart=/bin/sleep 5
