# tests/test_units/tests/partof.service
[Unit]
Description=PartOf service
PartOf=base.service

[Service]
ExecStart=/bin/sleep 100
