# tests/test_units/tests/poweroff.target
[Unit]
Description=test poweroff target
