# tests/test_units/tests/reliable_001.service
[Unit]
Description=reliable test 001

[Service]
ExecStart=/bin/sleep 888
