# tests/test_units/tests/reliable_002.service
[Unit]
Description=reliable test 002

[Service]
ExecStart=/bin/sleep 999

[Install]
RequiredBy=base.service
