# tests/test_units/tests/requires.service
[Unit]
Description=Requires service
Requires=base.service

[Service]
ExecStart=/bin/sleep 100
