# tests/test_units/tests/restart_001.service
[Unit]
Description=restart 001 service
JobTimeoutSec=2

[Service]
RestartPreventExitStatus=SIGPIPE
ExecStart=/bin/sleep 1
Restart=
