# tests/test_units/tests/wants.service
[Unit]
Description=Wants service
Wants=base.service

[Service]
ExecStart=/bin/sleep 100
