# use/devmaster替代udev运行/service/devmaster-simu-udev.service
[Unit]
DefaultDependencies=false

[Service]
Type=simple
RemainAfterExit=true
ExecStart=-/bin/bash /lib/devmaster/simulate_udev.sh
