# units/multi-user.target
[Unit]
Description=Multi user target
Requires=basic.target
