# fcty/usr/lib/sysmaster/system/sockets.target
[Unit]
Description=Sockets target
Documentation=man sysmaster special
