# fcty/usr/lib/sysmaster/system/sshd-keygen.target
[Unit]
Wants=sshd-keygen@rsa.service sshd-keygen@ecdsa.service sshd-keygen@ed25519.service
PartOf=sshd.service
