# fcty/usr/lib/sysmaster/system/sysinit.target
[Unit]
Description=system initialization target
Documentation=man sysmaster secial target
Wants=sysctl.service syslog.target
